# barlights.service
[Unit]
Description=Lights for the bar.

[Service]
ExecStart=/usr/sbin/lxdm
Restart=always
IgnoreSIGPIPE=no

[Install]
Alias=barlights.service

